# deploy/lepton-web.service
[Unit]
Description=Lepton Thermal Camera Web Interface
After=network.target
Wants=network.target

[Service]
Type=simple
User=@USER@
WorkingDirectory=@DIR@
ExecStart=@DIR@/.venv/bin/python @DIR@/lepton_web.py --host 0.0.0.0 --port @PORT@
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
